# site/multblog-supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /data/site/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
